# deployment/systemd/recruit-seeds-api.service
[Unit]
Description=Recruit Seeds API Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/recruit-seeds/apps/api
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=recruit-seeds-api

[Install]
WantedBy=multi-user.target