# deployment/systemd/recruit-seeds-jobs.service
[Unit]
Description=Recruit Seeds Jobs Application  
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/recruit-seeds/apps/jobs
Environment=NODE_ENV=production
Environment=PORT=3002
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=recruit-seeds-jobs

[Install]
WantedBy=multi-user.target